# systemd/user/mako.service
[Unit]
Description=A lightweight Wayland notification daemon 
Documentation=https://github.com/emersion/mako
BindsTo=sway-session.target

[Service]
Type=simple
ExecStart=/usr/bin/mako

[Install]
WantedBy=sway-session.target
